# Library/LaunchAgents/org.v2fly.v2ray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>org.v2fly.v2ray</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/v2ray/v2ray</string>
        <string>-config</string>
        <string>/usr/local/v2ray/config.json</string>
    </array>
    <key>RootDirectory</key>
    <string>/usr/local/v2ray</string>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
